# zram/zram0
#!/sbin/openrc-run

description="Configure zram0 as swap"

depend()
{
	provide zram0
	use dev-mount
	before dev
}

start()
{
	ebegin "Configure zram0 as swap"
	modprobe zram && zramctl zram0 -s 12G -t 6 -a lzo-rle && mkswap /dev/zram0 && swapon /dev/zram0 -p 100
	eend $?
}

stop()
{
	ebegin "Disabling zram0 as swap"
	swapoff /dev/zram0 && zramctl -r zram0
	eend $?
}
